# files/etc/systemd/system/sdr-scan.service
[Unit]
Description=SDR FM band scan
After=network-online.target
Conflicts=sdr-fm@active.service

[Service]
Type=oneshot
User=radio
TimeoutStartSec=600
ExecStartPre=-/usr/bin/sudo /usr/bin/systemctl stop sdr-fm@active
ExecStart=/usr/bin/python3 /opt/sdr-tuner/fm_scan.py --gain 40 --duration 60
ExecStartPost=-/usr/bin/sudo /usr/bin/systemctl start sdr-fm@active

[Install]
WantedBy=multi-user.target
