# files/etc/systemd/system/sdr-mux.service
[Unit]
Description=SDR FM-multistation supervisor (one pipeline per channel)
After=network-online.target icecast2.service sdr-iq-capture.service
Wants=network-online.target
# The mux needs the wideband IQ feed; starting it starts capture, which in turn
# Conflicts=sdr-fm@active (so legacy mono stops). Stopping the mux does NOT auto-
# stop capture — Flask does the explicit rollback (stop mux + capture, start mono).
Requires=sdr-iq-capture.service
Conflicts=sdr-fm@active.service

[Service]
Type=simple
User=radio
EnvironmentFile=/etc/sdr-streams/mux.env
RuntimeDirectory=sdr-streams
RuntimeDirectoryMode=0755
RuntimeDirectoryPreserve=yes
ExecStart=/opt/sdr-tuner/mux_supervisor.py
# SIGHUP = reload channels.json without dropping unchanged channels.
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=10
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
